# cronjob/cron.weekly/fellows.cron.file
#This cronjob will be executed at 4am every Monday
0 4 * * 1 www-data cron.sh

#Variable decleration

SHELL=/bin/bash
PATH=/sbin:/sbin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/
