# cronjob/cron.weekly/fellows.cron.file
#This cronjob will be executed at 4am every Monday
0 4 * * 1 root /cronjob/cron.sh

#Variable decleration

SHELL=/bin/bash
PATH=/sbin:/sbin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/
DATADIR=/Public/web/data
VISUAL_GIT=/tmp/Visualisation
