# cronjob/dashboard.cron
#This cronjob will be executed at 4am every Monday
0 4 * * 1 /bin/bash <path to the dashboard repository> cron.sh
